# profile.d/go.sh
#!/bin/sh

# use go modules
export GOPATH="$HOME/go"
export GO111MODULE=on
