# profile.d/log.sh
#!/bin/sh

# '~/.bzr.log' just won't stay deleted for some reason. Disable it altogether!
export BZR_LOG=/dev/null
